# .config/systemd/user/offlineimap.service
[Unit]
Description=Offlineimap
PartOf=graphical-session.target

[Install]
WantedBy=graphical-session.target

[Service]
ExecStart=offlineimap
